# crontab.sh.crontab.sample
$cron_unique_label
0	*	*	*	*	AUTOSSH_PORT=0 $PWD/autossh.sh

