# crontab.sh.crontab.sample
$cron_unique_label
0	*	*	*	*	$PWD/autossh.sh

